# pi-system/deploy/systemd/led-kiosk.service
[Unit]
Description=LED Touch Kiosk
After=graphical.target led-backend.service
Requires=led-backend.service

[Service]
Type=simple
User=ledvives
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/ledvives/.Xauthority
ExecStart=/usr/bin/chromium --kiosk --noerrdialogs --disable-infobars --touch-events=enabled --overscroll-history-navigation=0 --force-device-scale-factor=1.2 --user-data-dir=/tmp/chromium-kiosk --no-first-run --disable-dev-shm-usage http://127.0.0.1:3001
Restart=always
RestartSec=5

[Install]
WantedBy=graphical.target
